# deploy/shared/supervisor.app.conf
[program:tangent-django-boilerplate]
command=uwsgi --ini /opt/uwsgi.ini

[program:security-updates]
command=uwsgi --ini /opt/security-updates.ini
