# deploy/shared/supervisor.app.conf
[program:uwsgi-application]
command=uwsgi --ini /var/www/wsgi/uwsgi.ini

[program:uwsgi-security-updates]
command=uwsgi --ini /var/www/wsgi/security-updates.ini
